# crontab
*/15 * * * * python /app/commentFetcher.py
* * * * */7 /app/db_cleaner.sh
